# conf/supervisor/php.conf
[program:php-fpm]
command=php-fpm7.3 -F
autostart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
